# remarkable-fail.service
[Unit]
Description=reMarkable fail application

[Service]
ExecStart=/usr/sbin/remarkable-fail.sh
Restart=on-failure

